# control/memcached.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/memcached" ]; then
  invoke-rc.d memcached stop || exit $?
fi

